# source/planet-4-151612/wordpress/etc/service/php/run
#!/bin/bash
set -euf

source /app/bin/colours.sh
source /app/bin/env.sh

_good "Starting service: php-fpm${PHP_VERSION} ..."

# start PHP
exec php-fpm${PHP_VERSION} -c /etc/php/${PHP_VERSION}/fpm
